# deployment/relay.service
[Unit]
Description=Relay server (user)
After=network.target

[Service]
Type=simple
ExecStart=relay serve --addr :9292
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
